# com.jobhunt.briefing.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <!-- Daily JobHunt brief. Generated by /jobhunt-setup.
       Replace __PYTHON__, __MCPDIR__, and __HOME__ before installing to
       ~/Library/LaunchAgents/com.jobhunt.briefing.plist, then:
         launchctl load ~/Library/LaunchAgents/com.jobhunt.briefing.plist -->

  <key>Label</key>
  <string>com.jobhunt.briefing</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__MCPDIR__/daily_briefing.py</string>
  </array>

  <!-- Weekdays at 09:00 local time. launchd handles DST. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__MCPDIR__/briefing.log</string>
  <key>StandardErrorPath</key>
  <string>__MCPDIR__/briefing.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

</dict>
</plist>
